# packages/silk-deployment/silk-deployment-0.3.11.tar.gz/silk-deployment-0.3.11/silk/conf_templates/supervisord.conf
[program:%(site)s_%(dts)s]
command=%(cmd)s
directory=%(root)s
user=nobody
autostart=true
autorestart=true
stdout_logfile=%(root)s/logs/supervisor.log
redirect_stderr=True
environment=%(process_env)s
