# packages/baruwa/baruwa-1.1.1.tar.bz2/baruwa-1.1.1/extras/packaging/baruwa.cron.monthly
#!/bin/sh
#
# Baruwa
#
# Send PDF reports
baruwa-admin sendpdfreports >/dev/null 2>&1