# packages/baruwa/baruwa-1.1.1.tar.bz2/baruwa-1.1.1/extras/packaging/baruwa.cron
#!/bin/sh
#
# baruwa - 1.1.0
#

# send quarantine reports
# moved to baruwa.cron.d
# baruwa-admin sendquarantinereports >/dev/null

# clean quarantine 
baruwa-admin cleanquarantine >/dev/null 2>&1
# clean up the DB
baruwa-admin dbclean >/dev/null 2>&1
# update sa rule definitions
baruwa-admin updatesarules >/dev/null 2>&1
# clean up stale sessions
baruwa-admin cleanup >/dev/null 2>&1
# update geoip database
perl /usr/share/doc/GeoIP-*/fetch-geoipdata.pl
