# packages/OpenProximity/OpenProximity-12.04.1.tar.gz/OpenProximity-12.04.1/upstart/bluetoothd.conf
# bluetoothd - default configuration
#

description	"BlueZ startup"

start on starting network-manager
stop on stopping network-manager

respawn

script
    exec 1>>/var/log/openproximity/bluez.log
    exec 2>>/var/log/openproximity/bluez.err
    exec bluetoothd -n -d
    sleep 1
end script
