# packages/OpenProximity/OpenProximity-12.04.1.tar.gz/OpenProximity-12.04.1/upstart/openproximity-rpc-scanner.conf
# openproximity-scanner - default configuration
#

description	"OpenProximity Scanner"

start on starting openproximity-rpc
stop on stopping openproximity-rpc

respawn

export HOME=/root
export DEBUG_DEBUG_FILENAME=scanner
export OP_PATH=/opt/openproximity
export PATH=$OP_PATH/bin:$PATH

script
    cd /opt/openproximity
    exec 1>>/var/log/openproximity/scanner.log
    exec 2>>/var/log/openproximity/scanner.err
    exec bin/OpenProximity-manage rpcclient --rpc_client_client_mode=scanner
    sleep 30
end script
